# playbooks/files/jaeger-mcp.service.j2
[Unit]
Description=Jaeger MCP Bridge
After=network.target

[Service]
User=deploy
WorkingDirectory={{ jaeger_app_dir }}
Environment="JAEGER_URL={{ jaeger_url }}"
ExecStart={{ jaeger_app_dir }}/venv/bin/uvicorn app:app --host 0.0.0.0 --port {{ jaeger_port }}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target